# supervisord.conf
[supervisord]
nodaemon=true

[program:bot]
command=python bot.py
autostart=true
autorestart=true
stopsignal=SIGINT  # يمكنك تعيين إشارة إيقاف مختلفة إذا لزم الأمر
