# supervisord.conf
[supervisord]
nodaemon=true

[program:bot]
command=python bot.py
autostart=true
autorestart=true